# com.jack.daily-market-report-eod.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
      Daily Market Report — End-of-Day Scorecard
      Runs at 4:15 PM ET (16:15), Monday–Friday, after US markets close.
      Grades today's morning predictions against today's closing prices
      and regenerates report.html with the results.

      Install with:   ./install_schedule.sh
      Uninstall with: ./uninstall_schedule.sh
    -->
    <key>Label</key>
    <string>com.jack.daily-market-report-eod</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd "$HOME/Documents/Claude/Projects/Daily Report" &amp;&amp; ./venv/bin/python3 daily_market_report.py --eod --no-ai --no-open &amp;&amp; /usr/bin/open "$HOME/Documents/Claude/Projects/Daily Report/report.html"</string>
    </array>

    <!-- 4:15 PM ET, Mon (1) through Fri (5). One <dict> per day. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>15</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>15</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>15</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>15</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>15</integer>
        </dict>
    </array>

    <!-- Do not auto-run on load — only fire at the scheduled time. -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/daily-market-report-eod.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/daily-market-report-eod.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
